# chrooted/files/etc/systemd/system/restic.timer
[Unit]
Description=Timer for system backups

[Timer]
OnCalendar=hourly
Unit=restic.service

[Install]
WantedBy=timers.target
