# chrooted/files/etc/systemd/system/restic.service
[Unit]
Description=System backup

[Service]
ExecStart=/usr/local/bin/restic-backup
