# chrooted/etc/systemd/systemd/efi.mount
[Unit]
Description=efi partition mounting

[Mount]
What=%EFI_PART_PATH%
Where=/efi

[Install]
WantedBy=multi-user.target
